# kartoshka_bot.service
[Unit]
Description=Kartoshka Bot Telegram Service
After=network.target

[Service]
Type=simple
User=youruser
Group=youruser
WorkingDirectory=/home/youruser/kartoshka_bot
EnvironmentFile=/home/youruser/kartoshka_bot/.env
ExecStart=/home/youruser/kartoshka_bot/venv/bin/python3 /home/youruser/kartoshka_bot/kartoshka_bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
